# systemd/morinscode-cert-renew.service
[Unit]
Description=Renew morinscode.dev TLS certificates

[Service]
Type=oneshot
ExecStart=/usr/local/bin/renew-morinscode-certificates
