# systemd/morinscode-cert-renew.timer
[Unit]
Description=Run morinscode.dev certificate renewal twice daily

[Timer]
OnCalendar=*-*-* 03,15:00:00
Persistent=true
RandomizedDelaySec=30m

[Install]
WantedBy=timers.target
